# templates/default/service/upstart.conf.erb
# /etc/init/puma.conf - Puma config

description "Puma Application Server"

start on runlevel [2345]
stop on runlevel [06]

setuid <%= @user %>
setgid <%= @group %>

respawn
respawn limit <%= @respawn_count || 10 %> <%= @respawn_delay || 5 %>

script
# this script runs in /bin/sh by default
# respawn as bash so we can source in rbenv/rvm
# quoted heredoc to tell /bin/sh not to interpret
# variables

# source ENV variables manually as Upstart doesn't, eg:
#. /etc/environment

exec /bin/bash <<'EOT'
  # set HOME to the setuid user's home, there doesn't seem to be a better, portable way
  export HOME="$(eval echo ~$(id -un))"
  export RAILS_ENV="<%= @environment %>"
  export RACK_ENV="<%= @environment %>"

  source /etc/profile.d/rbenv.sh
  source /etc/profile.d/nvm.sh

  cd <%= @cwd %>
  logger -t puma "Starting server: <%= @name || "puma" %>"

  exec bundle exec puma -C config/puma.rb
EOT
end script
